# node_naturvinduet.conf
# node_naturvinduet.conf

description "generated by node-upstarter"

start on started networking
stop on runlevel [016]


limit nofile 1000000 1000000


console log

script
  
  mkdir -p /home/ubuntu/node_naturvinduet
  cd /home/ubuntu/node_naturvinduet
  
  node index.js
end script

respawn
